# meta-fsl-imx/recipes-bsp/pmic-setup-script/files/pmic_setup.service
[Unit]
Description=PMIC setup 

[Service]
Type=one shot
ExecStart=/tools/pmic_setup.sh
StandardOutput=kmsg+console

[Install]
WantedBy=multi-user.target
