# meta-tools/recipes-test/poweronoff-script/files/adv-power-test.service
[Unit]
Description=Advantech Power Test
After=systemd-logind

[Service]
Type=oneshot
ExecStart=/tools/power-test.sh

[Install]
WantedBy=multi-user.target
