# meta-tools/recipes-utils/add-router/files/adv-router.service
[Unit]
Description=Software Router Daemon
After=network.target

[Service]
ExecStart=/tools/adv-router

[Install]
WantedBy=multi-user.target
