# meta-tools/recipes-utils/add-virtual-network/files/add_virtual_network.service
[Unit]
Description=Add virtual network

[Service]
Type=idle
ExecStart=/tools/add_virtual_network.sh
StandardOutput=kmsg+console


[Install]
WantedBy=multi-user.target

