# meta-tools/recipes-utils/add-led-ctrl/files/adv-led-ctrl.service
[Unit]
Description=ADV LED ctrl Daemon
After=user@.service

[Service]
ExecStart=/tools/adv-led-ctrl.sh

[Install]
WantedBy=multi-user.target
