# meta-tools/recipes-utils/hwclock/files/hwclock.service
[Unit]
Description=hwclock Service

[Service]
Type=oneshot
ExecStart=/tools/hwclock.sh

[Install]
WantedBy=multi-user.target
