# meta-tools/recipes-test/advresize-script/files/adv_resize.service
[Unit]
Description=ADV Resize

[Service]
Type=idle
ExecStart=/tools/adv_resize.sh
StandardOutput=kmsg+console

[Install]
WantedBy=multi-user.target
