# meta-tools/recipes-utils/cellular-guard/files/cellular-guard.service
[Unit]
Description=Cellular Guard service
After=ModemManager.service

[Service]
ExecStart=/tools/cellular-guard.sh


[Install]
WantedBy=multi-user.target

